# crontab/root.cron
# matomo-archive 
# update need resatrt
* * * * *	/usr/local/bin/php /var/www/html/console core:archive --url=http://localhost:8080 > /home/www-data/matomo-archive.log